# contrib/systemd/tesla-charging-start.service
[Unit]
Description=Starts charging

[Service]
User=pi
Type=simple
# make sure to set correct paths
EnvironmentFile=/home/pi/bin/private/credentials
ExecStart=/home/pi/bin/tesla-api/run.py -i <vehicle_id> --charging-start
WorkingDirectory=/home/pi/bin/tesla-api/
